# templates/default/lol.conf.erb
#
# /etc/init/lol.conf
#
# Generated by Chef for <%= node[:fqdn] %>
#

description "LOL server"

start on runlevel [2345]
stop on runlevel [!2345]
setuid <%= @user %>
setgid <%= @group %>
respawn
respawn limit 3 30
env HOME=<%= @home %>
env PATH=<%= @home %>/.rbenv/shims:<%= @home %>/.rbenv/bin:/usr/local/bin:/usr/bin:/bin


pre-start script
    test -x <%= @ruby %> || { stop; exit 0; }
    test -e <%= @app_path %> || { stop; exit 0; }
    # this nasty hack would not be required if we used supervisor, or
    # if postgresql was upstart aware
    /etc/init.d/postgresql start
end script


chdir <%= @app_home %>
exec <%= @rerun %><%= @ruby %> <%= @app_path %>
