# templates/default/lol.conf.erb
#
# /etc/init/lol.conf
#
# Generated by Chef for <%= node[:fqdn] %>
#

description "LOL server"

start on stopped rc
stop on runlevel [016]
setuid <%= @user %>
setgid <%= @group %>
respawn
respawn limit 3 30
env HOME=<%= @home %>
env PATH=<%= @home %>/.rbenv/shims:<%= @home %>/.rbenv/bin:/usr/local/bin:/usr/bin:/bin


chdir <%= @app_home %>
exec <%= @rerun %><%= @ruby %> <%= @app_path %>
